# docker-dev/crontab
* * * * * /usr/local/bin/php /var/www/app/artisan schedule:run >> /var/log/app/schedule.log 2>&1
